# feature/bluetooth/systemd/smart-speaker-bluetooth.service
[Unit]
Description=SmartAudio Bluetooth Adapter Configuration
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=oneshot
ExecStart=/opt/bluetooth/bluetooth_manager.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
